# deploy/hdcf-control.supervisor.conf
[program:hdcf-control]
command=/usr/local/bin/hdcf-control -config /etc/hdcf/control-config.json
directory=/opt/hdcf
autorestart=true
autostart=true
startsecs=5
startretries=3
user=hdcf
stderr_logfile=/var/log/supervisor/hdcf-control.err.log
stdout_logfile=/var/log/supervisor/hdcf-control.out.log
environment=HDCF_CONTROL_CONFIG="/etc/hdcf/control-config.json"
stopsignal=TERM
stopasgroup=true
killasgroup=true
